# deploy/market-analysis.service
[Unit]
Description=Evidence-backed life-insurance market research
After=network-online.target business-analysis.service
Wants=network-online.target

[Service]
Type=oneshot
User=market-ai
Group=market-analysis
WorkingDirectory=/opt/business-analysis
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="HOME=/var/lib/business-analysis-market/home"
Environment="MARKET_ANALYSIS_DATA_DIR=/var/lib/business-analysis-market"
Environment="MARKET_ANALYSIS_INTERNAL_API=http://127.0.0.1:45679"
EnvironmentFile=/etc/business-analysis-market/market-analysis.env
ExecStart=/opt/business-analysis/backend/venv/bin/python /opt/business-analysis/backend/run_market_research.py
TimeoutStartSec=75min
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/business-analysis-market /var/log/business-analysis-market
UMask=0027
MemoryMax=2G
CPUQuota=180%

[Install]
WantedBy=multi-user.target
