# deploy/systemd.service
[Unit]
Description=Business Analysis Dashboard API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/business-analysis/backend
Environment="PATH=/opt/business-analysis/backend/venv/bin"
Environment="APP_ENV=production"
Environment="BUSINESS_ANALYSIS_DB=/opt/business-analysis/backend/business_data.db"
EnvironmentFile=-/opt/business-analysis/deploy/.admin_env
EnvironmentFile=-/opt/business-analysis/deploy/.ai_env
ExecStart=/opt/business-analysis/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 45679
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
